# fpms.service.template
[Unit]
Description=FPMS daemon responsible for OLED management

[Service]
WorkingDirectory=/usr/local/share/fpms
ExecStart=/usr/local/share/fpms/NanoHatOLED

[Install]
WantedBy=multi-user.target
